# dist/tpm-event-logging.service
[Unit]
Description=Set group tss for ima and system logging

[Service]
Type=oneshot
User=root
ExecStart=/bin/sh -c "/bin/chgrp tss /sys/kernel/security/tpm[0-9]/binary_bios_measurements /sys/kernel/security/ima/binary_runtime_measurements"

[Install]
WantedBy=multi-user.target
